# contrib/terraform-testing-infrastructure/modules/config-files/templates/hadoop-datanode.service.tpl
[Unit]
Description=Hadoop DataNode start/stop
After=remote-fs.target

[Service]
Environment=JAVA_HOME=${java_home}
Environment=HADOOP_HOME=${software_root}/hadoop/hadoop-${hadoop_version}
Environment=HADOOP_LOG_DIR=${hadoop_dir}/logs
User=hadoop
Group=hadoop
Type=oneshot
ExecStart=${software_root}/hadoop/hadoop-${hadoop_version}/bin/hdfs --daemon start datanode
ExecStop=${software_root}/hadoop/hadoop-${hadoop_version}/bin/hdfs --daemon stop datanode
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

